# docker/bacula-sd/bacula-9.4.4/platforms/systemd/bacula-dir.service.in
# Systemd Bacula service file
#
# Copyright (C) 2000-2017 Kern Sibbald
# License: BSD 2-Clause; see file LICENSE-FOSS
#
# /lib/systemd/system/bacula-dir.service
#
# Description:
#    Used to start/stop/reload the bacula director daemon service (bacula-dir)
#
#    enable : systemctl enable bacula-dir
#    start : systemctl start bacula-dir
#
#

# From http://www.freedesktop.org/software/systemd/man/systemd.unit.html
[Unit]
Description=Bacula Director Daemon service
Requires=network.target
After=network.target multi-user.target
RequiresMountsFor=@working_dir@ @sysconfdir@ @sbindir@

# From http://www.freedesktop.org/software/systemd/man/systemd.service.html
[Service]
Type=simple
User=@dir_user@
Group=@dir_group@
ExecStart=@sbindir@/bacula-dir -fP -c @sysconfdir@/bacula-dir.conf
ExecReload=@sbindir@/bacula-dir -t -c @sysconfdir@/bacula-dir.conf
ExecReload=/bin/kill -HUP $MAINPID
SuccessExitStatus=15

[Install]
WantedBy=multi-user.target
