# docker/baculum-web/confs/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/var/log/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
logfile_backups=10  
logfile_maxbytes=50MB  


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///var/run/supervisor.sock 


[program:php-fpm7.4]
command=/usr/sbin/php-fpm7.4 --nodaemonize -c /etc/php/7.4/fpm/php-fpm.conf
priority=50
autorestart=true


[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
priority=50
autorestart=true

